# familylink-addon/rootfs/etc/services.d/familylink/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# S6 service finish script for Family Link Auth
# ==============================================================================

bashio::log.info "Family Link Auth Service stopped"
